# crontab
*/3 * * * * /usr/local/bin/python /app/app/testNVR.py >> /app/cron.log 2>&1
